# python-app.conf
[program:pythonapp]
directory=/opt/app
command=/opt/app/env/bin/honcho start -f ./procfile run
autostart=true
autorestart=true
user=playlist-bot
# Environment variables ensure that the application runs inside of the
# configured virtualenv.
environment=VIRTUAL_ENV="/opt/app/env",PATH="/opt/app/env/bin",HOME="/opt/app/env/bin",USER="playlist-bot"
stdout_logfile=syslog
stderr_logfile=syslog
